# telegram-claude-yolo-bot.service
[Unit]
Description=Telegram Claude YOLO Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/telegram-claude-yolo-bot start
Restart=always
RestartSec=5
Environment=RUST_LOG=info

# Load environment variables from .env file
EnvironmentFile=-/etc/telegram-claude-yolo-bot/.env

[Install]
WantedBy=multi-user.target
